# Server/crontab
@hourly php ngWebPath/backend/at_home.php weather
*/7 * * * * php ngWebPath/backend/script/ddns.php