# Server/crontab
@hourly php ngWebPath/backend/at_home.php weather
